# wstunsrv.conf
# Ubuntu upstart config file for wstunsrv service
# Place in /etc/init/wstunsrv.conf

description "Websockets tunnel server"

start on runlevel [345]
stop on runlevel [!345]

env LOGFILE='/var/log/wstunsrv.log'

# automatically restart the service if it dies...
respawn
respawn limit 4 60

#oom score 100

# prevent the service from starting if the binary is missing
pre-start script
  echo checking wstunsrv pre-conditions
  # don't start if the binary doesn't exist
  [ ! -x /usr/local/bin/wstunsrv ] && { echo "/usr/local/bin/wstunsrv missing"; stop; exit 0; }

  # Ensure we can write to the log file
  touch "${LOGFILE}"
  chown www-data:www-data "${LOGFILE}"
  chmod 600 "${LOGFILE}"

  echo wstunsrv pre-conditions satisfied
end script

setuid www-data
setgid www-data

chdir /tmp
exec /usr/local/bin/wstunsrv -port 8000 -logfile ${LOGFILE}
